# app.cron
* * * * * /usr/local/bin/php /usr/src/app/bootstrap.php >> /var/log/cron/cron.log 2>&1
